# etc/init.d/sing-box-r
#!/sbin/openrc-run

name="sing-box-r"
description="Sing-box Transparent Proxy Service"

command="/usr/local/bin/sing-box"
command_args="run -C /etc/sing-box/conf/"
command_user="singbox"

pidfile="/run/sing-box.pid"
logfile="/var/log/sing-box.log"

depend() {
    need net
    after firewall
}

start_pre() {
    # 检查配置文件是否存在
    if [ ! -d "/etc/sing-box/conf" ]; then
        eerror "Configuration directory /etc/sing-box/conf not found!"
        return 1
    fi
    # 确保日志文件存在并有权限
    touch "$logfile"
    chown "$command_user" "$logfile"
}

start() {
    ebegin "Starting $name"
    # 添加 ip rule 和 ip route 规则
    ip rule add fwmark 0x1ed4 lookup 100 prio 32765 2>> "$logfile"
    if [ $? -eq 0 ]; then
        einfo "Added ip rule for fwmark 0x1ed4 to lookup table 100"
    else
        ewarn "Failed to add ip rule for fwmark 0x1ed4 or rule already exists"
    fi

    ip route add local default dev lo scope host table 100 2>> "$logfile"
    if [ $? -eq 0 ]; then
        einfo "Added ip route for table 100"
    else
        ewarn "Failed to add ip route for table 100 or route already exists"
    fi

    # 启动 sing-box
    start-stop-daemon --start --pidfile "$pidfile" \
        --exec "$command" --user "$command_user" \
        --background --make-pidfile \
        -- $command_args >> "$logfile" 2>&1
    eend $?

    # 在启动 sing-box 后应用 nftables 规则
    if [ -f "/etc/sing-box/r.nft" ]; then
        nft -f /etc/sing-box/r.nft && nft -f /etc/sing-box/cnlist.nft
        if [ $? -eq 0 ]; then
            einfo "Applied nftables rules from /etc/sing-box/r.nft"
        else
            eerror "Failed to apply nftables rules!"
        fi
    else
        ewarn "nftables rules file /etc/sing-box/r.nft not found!"
    fi
}

stop() {
    ebegin "Stopping $name"

    # Try stopping using the pidfile first
    start-stop-daemon --stop --quiet --pidfile "$pidfile"

    # If that fails, try killing the process by name
    if [ $? -ne 0 ]; then
        killall sing-box 2>> "$logfile"
        if [ $? -eq 0 ]; then
            einfo "Stopped sing-box process using killall"
        else
            ewarn "Failed to stop sing-box process"
            eend 1  # Indicate failure
            return 1
        fi
    fi

    eend 0  # Indicate success

    # 删除 ip rule 和 ip route 规则
    ip rule del fwmark 0x1ed4 lookup 100 prio 32765 2>> "$logfile"
    if [ $? -eq 0 ]; then
        einfo "Deleted ip rule for fwmark 0x1ed4"
    else
        ewarn "Failed to delete ip rule for fwmark 0x1ed4 or rule does not exist"
    fi

    ip route del local default dev lo scope host table 100 2>> "$logfile"
    if [ $? -eq 0 ]; then
        einfo "Deleted ip route for table 100"
    else
        ewarn "Failed to delete ip route for table 100 or route does not exist"
    fi
}

restart() {
    ebegin "Restarting $name"
    svc_stop
    svc_start
    eend $?
}
