# etc/init.d/sing-box
#!/sbin/openrc-run

name="sing-box"
description="Sing-box Transparent Proxy Service"

command="/usr/local/bin/sing-box"
command_args="run -C /etc/sing-box/conf/"
command_user="singbox"

pidfile="/run/sing-box.pid"
logfile="/var/log/sing-box.log"

depend() {
    need net
    after firewall
}

start_pre() {
    # 检查配置文件是否存在
    if [ ! -d "/etc/sing-box/conf" ]; then
        eerror "Configuration directory /etc/sing-box/conf not found!"
        return 1
    fi
    # 确保日志文件存在并有权限
    touch "$logfile"
    chown "$command_user" "$logfile"
}

start() {
    ebegin "Starting $name"



    # 启动 sing-box
    start-stop-daemon --start --pidfile "$pidfile" \
        --exec "$command" --user "$command_user" \
        --background --make-pidfile \
        -- $command_args >> "$logfile" 2>&1
    eend $?


}

stop() {
    ebegin "Stopping $name"

    # Try stopping using the pidfile first
    start-stop-daemon --stop --quiet --pidfile "$pidfile"

    # If that fails, try killing the process by name
    if [ $? -ne 0 ]; then
        killall sing-box 2>> "$logfile"
        if [ $? -eq 0 ]; then
            einfo "Stopped sing-box process using killall"
        else
            ewarn "Failed to stop sing-box process"
            eend 1  # Indicate failure
            return 1
        fi
    fi

    eend 0  # Indicate success

}

restart() {
    ebegin "Restarting $name"
    svc_stop
    svc_start
    eend $?
}
